# scripts/ops-monitor/systemd/ops-weekly-summary.timer
[Unit]
Description=Send Ops Monitor Weekly Summary

[Timer]
OnCalendar=Sat *-*-* 02:40:00
Persistent=true

[Install]
WantedBy=timers.target
